# debian/postinst
#! /bin/sh
#
update-rc.d modbusd defaults
service modbusd restart


