# debian/postrm
#! /bin/sh
#
update-rc.d modbusd remove
service modbusd stop

